# system/etc/systemd/system/user@.service.d/delegate.conf
# Docker rootless
[Service]
Delegate=cpu cpuset io memory pids
