# system/etc/systemd/system/ssh.service.d/00-socket.conf
[Unit]
After=ssh.socket
Requires=ssh.socket
