# system/etc/systemd/user/localsearch-3.service.d/backport-c968a2f43c03c5065d3c6b85670f777e7dd51057.conf
[Unit]
ConditionEnvironment=XDG_SESSION_CLASS=user
