# home-gpg/.config/systemd/user/gpg-agent.socket.d/home-gpg.conf
[Socket]
DirectoryMode=0310
ListenStream=
ListenStream=/dev/shm/home-gpg/S.gpg-agent
SocketMode=0660
