# system/etc/systemd/user/mpris-proxy.service.d/nonpipewire-disable.conf
# https://lore.kernel.org/linux-bluetooth/2488030.6sBWrkQq0l@daniel-desktop3/T/#u
[Install]
WantedBy=
WantedBy=pipewire.service
