# system/etc/systemd/system/systemd-resolved.service.d/backport-1ca315be0097bccc9ff55e09ac339a48fdb9a040.conf
[Unit]
After=var.mount
Requires=var.mount

[Service]
PrivateTmp=disconnected
