# system/etc/systemd/user/pipewire-pulse.service.d/nonpipewire-disable.conf
# https://gitlab.freedesktop.org/pipewire/pipewire/-/merge_requests/2600
[Unit]
BindsTo=pipewire.service
