# home/.config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/ibus-daemon
Type=Application
Name=ibus-daemon
StartupNotify=false
Terminal=false
