# system/etc/systemd/system/ssh.socket.d/addresses.conf
[Socket]
ListenStream=
ListenStream=2222
# F*** IPv6
BindIPv6Only=both
