# overlay/leader/system/etc/systemd/system/duckdns.timer
[Unit]
Description=Quarter-hour DuckDNS push

[Timer]
OnUnitInactiveSec=15m
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
